# swww/.config/systemd/user/change_wallpaper.timer
[Unit]
Description=Timer to change wallpapers every hour

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
